# crontab
# Run the renew process twice daily.
0 */12 * * * /renew.sh > /proc/1/fd/1 2>&1
